# templates/lighthouse.service.j2
[Unit]
Description=Nebula overlay networking tool - Lighthouse
Wants=basic.target network-online.target nss-lookup.target time-sync.target
After=basic.target network.target network-online.target
Before=sshd.service

[Service]
Type=simple
SyslogIdentifier=nebula
ExecReload=/bin/kill -HUP $MAINPID
ExecStart={{ nebula_install_destination }}/nebula -config {{ nebula_install_destination }}/config.yml
Restart=always
RestartSec=42s

[Install]
WantedBy=multi-user.target

